# scripts/file_checker/file_checker.service
[Unit]
Description=Check Flag File Service

[Service]
Type=oneshot
ExecStart=/test/file_checker.sh

[Install]
WantedBy=multi-user.target
