# scripts/file_checker/file_checker.timer
[Unit]
Description=Runs file_checker.service every minute

[Timer]
OnCalendar=*:0/1
Unit=file_checker.service

[Install]
WantedBy=timers.target